# supervisord.conf
[supervisord]
nodaemon=true

[program:deluged]
command = /usr/bin/deluged -d -c /config -L info -l /config/deluged.log

[program:deluge-web]
command = /usr/bin/deluge-web -c /config

